# test/fixtures/service/socket/foo.socket
# simple-socket.socket
[Unit]
Description=foo socket
PartOf=foo.service

[Socket]
ListenStream=127.0.0.1:8080

[Install]
WantedBy=sockets.target
